# curiefense/infuser-bq/infuser_bq.cron
* * * * * /scripts/bin/run.sh > /proc/1/fd/1 2>/proc/1/fd/2
# leave this line here